# swift_app/Resources/com.lune.openusagebar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.lune.openusagebar</string>
  <key>ProgramArguments</key><array>
    <string>__APP__/Contents/MacOS/OpenUsage Bar</string>
    <string>--background</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ProcessType</key><string>Interactive</string>
  <key>StandardOutPath</key><string>__HOME__/Library/Logs/OpenUsageBar.status.stdout.log</string>
  <key>StandardErrorPath</key><string>__HOME__/Library/Logs/OpenUsageBar.status.stderr.log</string>
</dict></plist>
